# dockerfiles/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/data/logs/supervisord.log
pidfile=/app/data/supervisord.pid

[program:redis]
command=redis-server /app/config/redis.conf --daemonize no
autostart=true
autorestart=true
priority=100
startsecs=5
startretries=3
stdout_logfile=/app/data/logs/redis.log
stderr_logfile=/app/data/logs/redis.log
# Resource limits for Redis
rlimit_as=1610612736  ; 1.5GB memory limit
rlimit_nofile=1024   ; File descriptor limit

[program:qdrant]
command=/usr/local/bin/qdrant
directory=/app
autostart=true
autorestart=true
priority=200
startsecs=15
startretries=3
stdout_logfile=/app/data/logs/qdrant.log
stderr_logfile=/app/data/logs/qdrant.log
environment=QDRANT__SERVICE__HTTP_PORT=6333,QDRANT__SERVICE__GRPC_PORT=6334,QDRANT__STORAGE__STORAGE_PATH=/app/data/qdrant,QDRANT__LOG_LEVEL=INFO
# Resource limits for Qdrant
rlimit_as=2147483648  ; 2GB memory limit
rlimit_nofile=2048   ; File descriptor limit

[program:neo4j]
command=/usr/bin/neo4j console
directory=/var/lib/neo4j
autostart=true
autorestart=true
priority=300
startsecs=30
startretries=3
stdout_logfile=/app/data/logs/neo4j.log
stderr_logfile=/app/data/logs/neo4j.log
# Run as app user (container runs as non-root)
# Resource limits for Neo4j
rlimit_as=3221225472  ; 3GB memory limit
rlimit_nofile=4096    ; File descriptor limit

[program:context-store]
command=python3 -m uvicorn src.mcp_server.main:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
priority=400
startsecs=10
startretries=3
stdout_logfile=/app/data/logs/context-store.log
stderr_logfile=/app/data/logs/context-store.log
# Resource limits for MCP Server
rlimit_as=1073741824  ; 1GB memory limit
rlimit_nofile=1024    ; File descriptor limit

# [program:health-monitor]
# command=/app/monitoring/health-monitor.sh
# directory=/app
# autostart=false
# autorestart=true
# priority=50
# stdout_logfile=/app/data/logs/health-monitor.log
# stderr_logfile=/app/data/logs/health-monitor.log
# Run health checks every 60 seconds
# startsecs=60
# startretries=3
# Temporarily disabled to focus on core services
